# torque/detect_nodes.conf
[program:init_torque]
command=/usr/local/bin/init.sh
startsecs=0
autorestart=false
