# slurm/detect_nodes.conf
[program:init_slurm]
command=/usr/local/bin/init.sh
startsecs=0
autorestart=false
